# supervisord.conf
[supervisord]
nodaemon=true
childlogdir=/var/log/supervisord

[program:xvfb]
command=Xvfb :99 -screen 0 1024x768x16
